# deploy/recipe-rotation.service
[Unit]
Description=Recipe rotation HTTP server
After=network.target

[Service]
Type=simple
Environment=PORT=8080
Environment=RECIPE_DB_PATH=/var/lib/recipe-rotation/recipes.db
WorkingDirectory=/var/lib/recipe-rotation
ExecStart=/usr/local/bin/recipe-rotation
Restart=on-failure
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
